# etc/skel/.config/autostart/icons.desktop
[Desktop Entry]
Type=Application
Exec=konsole -e /bin/sh -c "gnome-terminal -- bash -c "~/.config/winux7/auto/icon.sh; exec bash" & /usr/bin/sleep 1"
X-GNOME-Autostart-enabled=true
NoDisplay=false
Hidden=false
Name[C]=icons
Comment[C]=No description
X-GNOME-Autostart-Delay=1
